# sa-train.target
[Unit]
Description=Spamassassin learn Timer Target
StopWhenUnneeded=yes
